# recipes-openxt/xen-vif-scripts/xen-vif-scripts-ndvm/xen-vif-backend.rules
SUBSYSTEM=="xen-backend", KERNEL=="vif*|vwif*", RUN+="/etc/xen/scripts/vif $env{ACTION}"
