# recipes/linux-input/files/99-wacom-serial.rules

ACTION=="add|change",  SUBSYSTEM=="tty", ATTRS{id}=="WACf*", RUN+="/usr/sbin/inputattach --daemon --wacom-384 /dev/$name"

