# recipes-openxt/libargo/files/13-argo.rules
KERNEL=="argo_dgram", RUN+="/sbin/restorecon /dev/argo_dgram"
KERNEL=="argo_stream", RUN+="/sbin/restorecon /dev/argo_stream"
