# recipes-openxt/xen-tap-scripts/files/xen-tap-backend.rules
SUBSYSTEM=="xen-backend", KERNEL=="tap*", RUN+="/etc/xen/scripts/tap $env{ACTION}"
