# recipes-openxt/xen-legacy-block-scripts/files/xen-block-backend.rules
SUBSYSTEM=="xen-backend", KERNEL=="vbd*", RUN+="/etc/xen/scripts/block $env{ACTION}"
