# recipes-openxt/libv4v/files/13-v4v.rules
KERNEL=="v4v_dgram", RUN+="/sbin/restorecon /dev/v4v_dgram"
KERNEL=="v4v_stream", RUN+="/sbin/restorecon /dev/v4v_stream"
