# recipes-extended/xen/files/xen-vbd3-backend.rules
SUBSYSTEM=="xen-backend", KERNEL=="vbd3*", RUN+="/etc/xen/scripts/block-vbd3 $env{ACTION}"
